# systemd/system/debug-trigger@.service
[Unit]
Description=BMC Debug Trigger

[Service]
Type=simple
ExecStart=/usr/bin/debug-trigger /%I /proc/sysrq-trigger
